# usr/share/doc/va.script/examples/etc/cron.monthly/clean
#!/bin/bash

logger "$0"

va.disk.file.old.delete.sh /var/backups 45
