# usr/share/doc/va.script/examples/etc/cron.weekly/backup
#!/bin/bash

logger "$0"

va.bak.tar.sh /var/www /var/backups/ z quiet
