# usr/share/doc/va.script/examples/etc/cron.weekly/letsencrypt
#!/bin/bash

logger "$0"

certbot -q renew
