# usr/share/doc/va.script/examples/etc/cron.daily/backup
#!/bin/bash

logger "$0"

va.mysql.db.export.all.sh /var/backups z quiet
